# airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear -l /usr/local/bin/acreetionos-server-display %I 38400 linux
